# infrastructure/kodiak-daily-restart.service
[Unit]
Description=Restart kodiak
 
[Service]
ExecStart=/root/kodiak-restart.sh
